# container/etc/cron.d/burp_reports
# burp reports
0  9  * * 1 root     /usr/bin/burp-reports -c /etc/burp/burp-reports.conf --report email_outdated
0 11  * * 3 root     /usr/bin/burp-reports -c /etc/burp/burp-reports.conf --report email_outdated
